# templates/default/ldm.erb
#!/bin/bash
# Automatically Generated by chef, manual changes will be over-written
# MM20131204 adding clean to start
# MM20131010 first go at init script for ldm
###
# chkconfig: 2345 99 1
# description: Manages UniData LDM server
###

LDMUSER=<%= @user %>
LDMADM=<%= @ldmadmin %>

case "$1" in
	start)
		echo "${1}ing ldm "
		if [[ -f <%= @pidfile =%> && $(ps -ef |grep -v grep |grep -c $(cat <%= @pidfile =%>)) = 0 ]] ; then
			/bin/su - $LDMUSER -c "$LDMADM clean"
		fi
		/bin/su - $LDMUSER -c "$LDMADM $1"
		;;
	stop|restart)
		echo "${1}ing ldm "
		/bin/su - $LDMUSER -c "$LDMADM $1"
		;;
	status)
		echo -n "checking status of ldm: "
		/bin/su - $LDMUSER -c "$LDMADM isrunning" && echo "running" || echo "not running"
		;;
	*)
		echo "Usage: $0 <start|stop|restart|status>"
		exit 1
		;;
esac

exit 0

